# control.service
[Unit]
Description=Control
 
[Service]
ExecStart=/root/newcontrol.sh
 
[Install]
WantedBy=multi-user.target
