# kobo/etc/udev/rules.d/control.rules
KERNEL=="loop0", ACTION=="add", RUN+="/usr/local/control.sh"
